# etc/systemd/system/kiosk-config.service
[Unit]
Description=Apply per-device kiosk config from /data
Before=kiosk.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c '[ -f /data/poly-kiosk/config ] && cp -f /data/poly-kiosk/config /etc/default/c60-kiosk || true'

[Install]
WantedBy=multi-user.target
